# DEBIAN/postinst
systemctl enable projetor-on.service
systemctl enable projetor-off.service
systemctl start projetor-off.service