# etc/systemd/system/projetor-on.service
[Unit]
Description=Executado ao Ligar

[Service]
Type=oneshot
User=root
ExecStart=/usr/projector-lamp -on

[Install]
WantedBy=multi-user.target
